# files/default/teamcity-server.conf
# teamcity-server - TeamCity Server
# THIS FILE IS MANAGED BY CHEF. ALL CHANGES WILL BE OVERWRITTEN!
#

description "TeamCity Server"

start on runlevel [2345]
stop on runlevel [016]

# TeamCity can take a while to shutdown
kill timeout 30

respawn

pre-start script
    export TEAMCITY_SERVER_MEM_OPTS="-Xss2m -XX:SoftRefLRUPolicyMSPerMB=10000 -XX:+UseConcMarkSweepGC."
    exec /opt/teamcity/current/TeamCity/bin/teamcity-server.sh start
end script
post-stop exec /opt/teamcity/current/TeamCity/bin/teamcity-server.sh stop